# services/updater.service
[Unit]
Description=Comprueba e instala actualizaciones para fire-truck-app
After=network-online.target
Wants=network-online.target
# NO bloqueamos el arranque del sistema.


[Service]
Type=oneshot
User=cosigein
Group=cosigein
TimeoutStartSec=5min
ExecStart=/usr/bin/python3 /home/cosigein/fire-truck-app/setup.py

[Install]
WantedBy=multi-user.target